# debian/etc/inittab
T0:2345:respawn:/sbin/getty -L ttyS2 115200 vt100
